# backend/manebid-backend.service
[Unit]
Description=ManeBid — Convention Auction Manager Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=manebid
Group=manebid
WorkingDirectory=/opt/manebid
Environment=NODE_ENV=production
EnvironmentFile=/etc/manebid/manebid.env
ExecStart=/usr/bin/node /opt/manebid/backend.js
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
#SystemCallFilter=@system-service
CapabilityBoundingSet=
AmbientCapabilities=

ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
PrivateMounts=yes
RestrictNamespaces=yes


# Allow writes only where needed
ReadWritePaths=/var/lib/manebid /var/log/manebid
ReadOnlyPaths=/opt/manebid

[Install]
WantedBy=multi-user.target
